# etc/systemd/system/tiny-dfr.service.d/network.conf
# Drop-in that grants the tiny-dfr unit network access. Required for the weather
# buttons, which fetch from wttr.in over HTTPS. RestrictAddressFamilies merges
# across drop-ins, so this adds AF_INET/AF_INET6 to whatever the main unit sets
# (typically AF_UNIX AF_NETLINK) without overriding a customized service file.
[Service]
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6
